# bt_joy/debian/bt-joy-server.service
[Unit]
Description=bt-joy Betaflight MSP server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/joy-server --config /etc/bt-joy/server.yaml
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
